# deployment/supervisor-queue.conf
; Supervisor config for Sinemani Laravel queue worker.
; Install:
;   sudo cp deployment/supervisor-queue.conf /etc/supervisor/conf.d/sinemani-queue.conf
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start sinemani-queue:*
;
; Logs: /var/log/supervisor/sinemani-queue.log

[program:sinemani-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/dramabox/artisan queue:work --queue=default,notifications,video --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/sinemani-queue.log
stopwaitsecs=3600
